# gunicorn.ini
[program:gunicorn]
command=/vagrant/src/.venv/bin/gunicorn --bind 0.0.0.0 budgetapp.wsgi
directory=/vagrant/src/budgetapp
user=vagrant
autostart=true
autorestart=true
stderr_logfile=/var/log/gunicorn.err.log
stdout_logfile=/var/log/gunicorn.out.log
